# systemd/sub-manager.service
[Unit]
Description=Multi-Server Sub Manager Service
After=network.target

[Service]
User=root
WorkingDirectory=/opt/sub-manager
Environment="PROJECT_DIR=/opt/sub-manager"
Environment="APP_PORT=666"
Environment="CACHE_TTL=30"
Environment="WEB_PATH=my-panel"
Environment="ALLOW_ORIGINS=http://localhost:5173,http://127.0.0.1:5173"
Environment="VERIFY_TLS=true"
Environment="CA_BUNDLE_PATH="
Environment="READ_ONLY_MODE=false"
Environment="SUB_RATE_LIMIT_COUNT=30"
Environment="SUB_RATE_LIMIT_WINDOW_SEC=60"
ExecStartPre=/usr/bin/bash -c "/usr/bin/fuser -k 666/tcp || true"
ExecStart=/opt/sub-manager/venv/bin/uvicorn main:app --host 127.0.0.1 --port 666
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
