# systemd/sub-manager.service
[Unit]
Description=Multi-Server Sub Manager Service
After=network.target

[Service]
User=root
WorkingDirectory=/opt/sub-manager
Environment="PROJECT_DIR=/opt/sub-manager"
Environment="APP_PORT=666"
Environment="CACHE_TTL=30"
Environment="WEB_PATH=my-panel"
Environment="ALLOW_ORIGINS=http://localhost:5173,http://127.0.0.1:5173"
Environment="VERIFY_TLS=true"
Environment="CA_BUNDLE_PATH="
Environment="READ_ONLY_MODE=false"
Environment="SUB_RATE_LIMIT_COUNT=30"
Environment="SUB_RATE_LIMIT_WINDOW_SEC=60"
Environment="TRAFFIC_STATS_CACHE_TTL=10"
Environment="ONLINE_CLIENTS_CACHE_TTL=10"
Environment="TRAFFIC_STATS_STALE_TTL=120"
Environment="ONLINE_CLIENTS_STALE_TTL=60"
Environment="CLIENTS_CACHE_TTL=20"
Environment="CLIENTS_CACHE_STALE_TTL=180"
Environment="TRAFFIC_MAX_WORKERS=8"
Environment="COLLECTOR_BASE_INTERVAL_SEC=5"
Environment="COLLECTOR_MAX_INTERVAL_SEC=60"
Environment="COLLECTOR_MAX_PARALLEL=8"
Environment="REDIS_URL="
Environment="AUDIT_QUEUE_BATCH_SIZE=200"
Environment="ROLE_VIEWERS="
Environment="ROLE_OPERATORS="
Environment="MFA_TOTP_ENABLED=false"
Environment="MFA_TOTP_USERS="
Environment="MFA_TOTP_WS_STRICT=false"
ExecStartPre=/usr/bin/bash -c "/usr/bin/fuser -k 666/tcp || true"
ExecStart=/opt/sub-manager/venv/bin/uvicorn main:app --host 127.0.0.1 --port 666
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
